# webhook-deploy.service
[Unit]
Description=MC LXD Manager - GitHub Webhook Deployment Service
After=network.target

[Service]
Type=simple
User=cordt
WorkingDirectory=/home/cordt/repos/mclauncher
Environment="NODE_ENV=production"
Environment="WEBHOOK_PORT=9000"
Environment="WEBHOOK_SECRET=CHANGE_THIS_SECRET"
Environment="REPO_PATH=/home/cordt/repos/mclauncher"
Environment="DEPLOY_BRANCH=feature/advanced-mod-management"
ExecStart=/usr/bin/node /home/cordt/repos/mclauncher/webhook-deploy.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=webhook-deploy

[Install]
WantedBy=multi-user.target
