# deploy/mc-lxd-manager.service
[Unit]
Description=MC LXD Manager (Node)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/mc-lxd-manager.env
WorkingDirectory=/opt/mc-lxd-manager
# Ensure lxc is on PATH for this unit (snap installs typically place it in /usr/bin)
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/node apps/server/dist/index.js
Restart=always
RestartSec=3
User=mc
Group=mc
# Let mc use the lxd unix socket
SupplementaryGroups=lxd

[Install]
WantedBy=multi-user.target
